# rootfs/etc/s6-overlay/s6-rc.d/oneshot-update-nvim-configuration/up
#!/command/with-contenv execlineb

foreground {
    chmod +x /usr/local/bin/nvim-config-update
}

/usr/local/bin/nvim-config-update
